# scripts/com.novakai.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.novakai.watchdog</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/christopherdasca/Programming/Novakai-Command/scripts/nvk-watchdog.mjs</string>
    <string>watch</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/christopherdasca/Programming/Novakai-Command</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/christopherdasca/Programming/Novakai-Command/.novakai-command/watchdog.out</string>
  <key>StandardErrorPath</key><string>/Users/christopherdasca/Programming/Novakai-Command/.novakai-command/watchdog.out</string>
</dict>
</plist>
